# guifi-proxy3-1.1.24/debian/postrm
#! /bin/sh
# postrm script for guifi-proxy3
#
# see: dh_installdeb(1)

set -e

case "$1" in
        purge|remove)

                if [ -e /usr/share/debconf/confmodule ] ; then

                        . /usr/share/debconf/confmodule
                        db_purge

                fi

                if [ -x /usr/bin/ucf ] ; then 

                        /usr/bin/ucf --purge /etc/guifi-proxy3/config.sh
                        /bin/rm -rf /etc/guifi-proxy3/
                fi

		if [ -e /etc/squid3/squid.conf.backup.initial ]; then

			echo "Restoring the initial configuration squid.conf.backup.initial ...... .......";
			FILE="/etc/squid3/squid.conf"
  			grep -Ev 'guifi' "$FILE" > "$FILE.tmp"
                        grep -Ev 'custom' "$FILE" > "$FILE.tmp"
			mv "$FILE.tmp" $FILE
		fi
		echo "Stoping the server squid ...... .......";
		        invoke-rc.d squid3 stop

                ;;

        upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
                ;;

        *)
                echo "prerm called with unknown argument \`$1'" >&2
                exit 1

esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

