# guifi-proxy3-0.0.1/debian/guifi-proxy3.cron.d
#
# Regular cron jobs for the guifi-proxy3 package
#
55 *	* * *	root	/usr/share/guifi-proxy/guifi-proxy.sh >> /var/log/guifi-proxy/guifi-proxy.log 2>&1;
